# wibed-system/files/etc/uci-defaults/rm-exit0-rclocal.sh
#!/bin/sh
sed s/"^exit 0"//g -i /etc/rc.local
