# wibed-system/files/etc/uci-defaults/try_prepare_usb.sh
#!/bin/sh

for device in sda sdb sdc; do
    if [ -b /dev/$device ]; then
        continue
    fi
done
