# wibed-system/files/etc/uci-defaults/try_prepare_usb.sh
#!/bin/sh

for device in sdb sdc sde sdf; do
    [ -b /dev/$device ] || continue
done
