# wibed-system/files/etc/uci-defaults/wdr-leds.sh
echo 'echo timer > /sys/class/leds/tp-link\:blue\:qss/trigger' >> /etc/rc.local
