# wibed-system/files/etc/uci-defaults/wibed-config.sh
#!/bin/sh

wibed-config 2>&1 > /root/wibed-config.log
