# wibed-system/files/etc/hotplug.d/button/10-wbm-restore
#!/bin/sh
if [ "$BUTTON" = wps ] && [ "$ACTION" = released ]
	then 
	echo "Restoring defaults..." 
	/sbin/jffs2reset -y 
	/sbin/reboot -f
fi
